# GTM/etc/xinetd.d/vista-cia
service vista-cia
{
  port = 9100
  socket_type = stream
  protocol = tcp
  type = UNLISTED
  user = foia
  server = /home/foia/bin/cia.sh
  wait = no
  disable = no
  per_source = UNLIMITED
  instances = UNLIMITED
}
